# systemd/bedrock-fetch-fundamentals.service
# Auto-generert av `bedrock systemd generate`.
# Ikke rediger manuelt — kjør CLI-kommandoen på nytt hvis config endrer seg.
[Unit]
Description=Bedrock fetch: fundamentals (bedrock.fetch.fred)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pc/bedrock
ExecStart=/home/pc/bedrock/.venv/bin/bedrock fetch run fundamentals
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
