# systemd/bedrock-fetch-cot_legacy.service
# Auto-generert av `bedrock systemd generate`.
# Ikke rediger manuelt — kjør CLI-kommandoen på nytt hvis config endrer seg.
[Unit]
Description=Bedrock fetch: cot_legacy (bedrock.fetch.cot_cftc)
After=network-online.target
Wants=network-online.target
OnFailure=bedrock-notify@%N.service

[Service]
Type=oneshot
WorkingDirectory=/home/pc/bedrock
ExecStart=/home/pc/bedrock/.venv/bin/bedrock fetch run cot_legacy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
