# systemd/bedrock-fetch-fundamentals.service
# Auto-generert av `bedrock systemd generate`.
# Ikke rediger manuelt — kjør CLI-kommandoen på nytt hvis config endrer seg.
[Unit]
Description=Bedrock fetch: fundamentals (bedrock.fetch.fred)
After=network-online.target
Wants=network-online.target
OnFailure=bedrock-notify@%N.service

[Service]
Type=oneshot
WorkingDirectory=/home/pc/bedrock
ExecStartPre=/bin/sh -c 'for _ in $(seq 1 30); do getent ahosts cloudflare.com >/dev/null 2>&1 && exit 0; sleep 2; done; exit 0'
ExecStart=/home/pc/bedrock/.venv/bin/bedrock fetch run fundamentals --fail-tolerance-pct 50
ExecStartPost=/home/pc/bedrock/scripts/signals_bot_regen.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
