# systemd/bedrock-fetch-prices.service
# Auto-generert av `bedrock systemd generate`.
# Ikke rediger manuelt — kjør CLI-kommandoen på nytt hvis config endrer seg.
[Unit]
Description=Bedrock fetch: prices (bedrock.fetch.prices)
After=network-online.target
Wants=network-online.target
OnFailure=bedrock-notify@%N.service

[Service]
Type=oneshot
WorkingDirectory=/home/pc/bedrock
ExecStart=/home/pc/bedrock/.venv/bin/bedrock fetch run prices
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
